# pokevend.service
[Unit]
Description=PokéVend — Pokemon card vending machine

[Service]
Type=simple
WorkingDirectory=__WORKING_DIR__
Environment=SDL_AUDIODRIVER=dummy
ExecStart=__WORKING_DIR__/venv/bin/python -m pokevend.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
